# com.jeffbai.ibkr-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jeffbai.ibkr-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jeffbai/repos/ibkr-sync/venv/bin/python</string>
        <string>/Users/jeffbai/repos/ibkr-sync/sync_portfolio.py</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>WorkingDirectory</key>
    <string>/Users/jeffbai/repos/ibkr-sync</string>

    <key>StandardOutPath</key>
    <string>/tmp/ibkr-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/ibkr-sync.err</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
